# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=100MB
logfile_backups=3
pidfile=/var/run/supervisord.pid
loglevel=info

[program:fhir-module]
command=gunicorn --preload -c gunicorn_config.py -w 1 -b 0.0.0.0:5000 main:app
directory=/opt/fhir-module
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1
environment=PROMETHEUS_MULTIPROC_DIR="/tmp/prometheus_multiproc"

[program:nextjs-ui]
command=node server.js
directory=/app
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2
environment=PORT="3000",HOSTNAME="0.0.0.0"


